# systemd/dronesim.target
[Unit]
Description=Drone simulation services

[Install]
WantedBy=default.target
