# systemd/liftoff-mavlink-bridge.service
[Unit]
Description=Liftoff MAVLink bridge for GCS
After=network.target zenohd.service

[Service]
Type=simple
ExecStart=%h/.cargo/bin/liftoff-mavlink-bridge $LIFTOFF_MAVLINK_BRIDGE_ARGS
Environment=RUST_LOG=info
EnvironmentFile=-%h/.config/liftoff/env
Restart=on-failure
RestartSec=3

[Install]
WantedBy=liftoff.target
