# systemd/liftoff.target
[Unit]
Description=Liftoff services (manual)

[Install]
WantedBy=default.target
